# deb_package/etc/skel/.config/autostart/libreoffice-loader.desktop
[Desktop Entry]
Exec=/usr/share/soffice-loader/libreup
Name=Libreoffice Loader
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
